# server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/logs/supervisord.log
logfile_maxbytes=100MB
logfile_backups=0

[program:node]
command=yarn start
directory=/usr/src/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx "-g daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
